# crontab
1 0 * * * root python /workday2mqtt/main.py > /var/log/cron.log 2>&1
